# ansible/roles/gse/files/gse.service
[Unit]
Description=GSE Flask daemon
After=network.target

[Service]
ExecStart=/usr/bin/gse.sh

[Install]
WantedBy=multi-user.target